# config/systemd/hyprwhspr-resume.service
[Unit]
Description=hyprwhspr suspend/resume handler
Before=sleep.target
StopWhenUnneeded=yes

[Service]
Type=oneshot
RemainAfterExit=yes
# On suspend: create marker file to signal suspend event
ExecStart=/bin/sh -c 'mkdir -p ~/.config/hyprwhspr && touch ~/.config/hyprwhspr/.suspend_marker'
# On resume: remove marker (main service detects absence = resume)
ExecStop=/bin/sh -c 'rm -f ~/.config/hyprwhspr/.suspend_marker'

[Install]
WantedBy=sleep.target
