# config/systemd/hyprwhspr.service
[Unit]
Description=hyprwhspr tts
Documentation=https://github.com/goodroot/hyprwhspr
After=graphical-session.target pipewire.service ydotool.service
Requires=ydotool.service
Wants=pipewire.service

[Service]
Type=simple
# Wait for Wayland environment - wl-copy needs WAYLAND_DISPLAY to be set
ExecStartPre=/bin/bash -c 'for i in $(seq 1 30); do [ -n "$WAYLAND_DISPLAY" ] && exit 0; sleep 0.5; done; echo "Warning: WAYLAND_DISPLAY not set after 15s"'
ExecStart=/usr/lib/hyprwhspr/bin/hyprwhspr
Environment=HYPRWHSPR_ROOT=/usr/lib/hyprwhspr
Restart=on-failure
RestartSec=2
StartLimitIntervalSec=10
StartLimitBurst=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target