# config/systemd/meeting-recorder.service
[Unit]
Description=meeting-recorder (hyprwhspr)
Documentation=https://github.com/goodroot/hyprwhspr

PartOf=graphical-session.target
After=graphical-session.target

After=pipewire.service
Wants=pipewire.service

[Service]
Type=simple
ExecStart=/usr/lib/hyprwhspr/bin/meeting-recorder
Environment=HYPRWHSPR_ROOT=/usr/lib/hyprwhspr
Environment=PYTHONUNBUFFERED=1
Restart=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
